# deploy/okx-robot.service
[Unit]
Description=OKX Copy Trading Robot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/okx-robot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/okx-robot
EnvironmentFile=/opt/okx-robot/.env
ExecStartPre=/usr/bin/test -f /opt/okx-robot/.env
ExecStartPre=/usr/bin/test -f /opt/okx-robot/config.yaml
ExecStart=/opt/okx-robot/.venv/bin/python -m src.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
